# templates/etc/cron.daily/gpgrefresh
#!/bin/sh
#
# Refresh the GPG keys on root's keyring
#

gpg --homedir /root/.gnupg \
    --keyserver <%= @keyserver %> \
    --refresh-keys 2> /dev/null
