# templates/etc/cron.daily/puppet_clean_reports
#!/bin/sh
find /var/lib/puppet/reports -type f -mtime +<%= @reports_days -%> -exec rm {} \+
