# templates/etc/cron.d/ldap
# Dump the LDAP databases every 5 minutes
*/10 * * * * root /usr/local/sbin/backup-ldap-database.sh
