# files/etc/network/if-up.d/ip6tables
#!/bin/sh
#
# Simple script to restore IPv6 iptables rules

/sbin/ip6tables-restore < /etc/ip6tables.conf

