# files/etc/network/if-up.d/iptables
#!/bin/sh
#
# Simple script to restore iptables rules

/sbin/iptables-restore < /etc/iptables.conf

