# templates/etc/cron.hourly/glpi
#!/bin/sh
sh /usr/share/glpi/plugins/massocsimport/scripts/ocsng_fullsync.sh
